# WLEDMatrixManager/rootfs/etc/services.d/example/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Add-on Service Start
# ==============================================================================
# FastAPI backend server that communicates with Home Assistant
# https://developers.home-assistant.io/docs/add-ons/
# ==============================================================================

# Declare variables
declare log_level

## Get log level from config
log_level=$(bashio::config 'log_level')

## Set Python environment
export PYTHONUNBUFFERED=1
export LOG_LEVEL="${log_level:=info}"

bashio::log.info "Starting Home Assistant Add-on..."
bashio::log.info "Log level: ${LOG_LEVEL}"

## Install custom integration into Home Assistant
if [ -d /homeassistant ]; then
    bashio::log.info "Installing WLED Matrix Manager integration..."
    mkdir -p /homeassistant/custom_components/wled_matrix_manager
    cp -rf /app/custom_components/wled_matrix_manager/* \
        /homeassistant/custom_components/wled_matrix_manager/
    bashio::log.info "Integration installed to /homeassistant/custom_components/"
else
    bashio::log.warning "HA config dir not found — integration NOT installed"
fi

## Change to app directory
cd /app/backend

## Run FastAPI application
exec python3 main.py
